# deploy/umstutorial-deploy.service
[Unit]
Description=Build and publish umstutorial to its webroot
After=network-online.target
Wants=network-online.target

[Service]
Type=oneshot
User=__DEPLOY_USER__
Group=__DEPLOY_USER__
WorkingDirectory=/srv/umstutorial
# Explicit PATH so deploy.sh finds node/npm/git/rsync. Override here if
# node lives outside the standard system locations (e.g. under nvm:
# prepend /home/<user>/.nvm/versions/node/v22.21.0/bin to this list).
Environment=PATH=/usr/local/sbin:/usr/local/bin:/usr/sbin:/usr/bin:/sbin:/bin
ExecStart=/srv/umstutorial/deploy/deploy.sh
Nice=10
TimeoutStartSec=10min
